# config/supervisord.conf
; ── QA Fleet supervisord config ───────────────────────────────────────────────
; NOTE: entrypoint.sh overwrites this file at container startup based on the
; project config (qa-fleet.conf). This file is a baked-in fallback only and
; is not used in normal operation.

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
